# gandiDns/rootfs/etc/services.d/updater/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: GandiDns
# Runs Updater script
# ==============================================================================


bashio::log.info "Starting Ip Updater..."

exec /usr/bin/updater.sh
